# deploy/docker/jos-daily-brief.timer
[Unit]
Description=Print the Jos Daily Brief every morning at 06:30

[Timer]
OnCalendar=*-*-* 06:30:00
Persistent=true

[Install]
WantedBy=timers.target
